# infra/systemd/pact-journal.service
[Unit]
Description=pact-journal — distributed immutable configuration log (Raft quorum)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pact-journal \
    --node-id ${PACT_JOURNAL_NODE_ID} \
    --listen ${PACT_JOURNAL_LISTEN} \
    --data-dir ${PACT_JOURNAL_DATA_DIR}
Restart=on-failure
RestartSec=5

EnvironmentFile=/etc/pact/journal.env

User=pact
Group=pact

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/pact/journal
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
