# infra/systemd/pact-agent.service
[Unit]
Description=pact-agent — per-node configuration management daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=+/bin/bash -c 'mkdir -p /var/lib/pact /run/pact && chown pact:pact /var/lib/pact /run/pact'
ExecStart=/usr/local/bin/pact-agent --config /etc/pact/agent.toml
Restart=on-failure
RestartSec=5

EnvironmentFile=/etc/pact/agent.env

User=pact
Group=pact

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/pact /run/pact
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
